# init/pmonitor.conf
# system monit 

start on stopped runlevel=[2345]
stop on runlevel [!$RUNLEVEL]

respawn
exec  /letv/script/process_monitor/pmanager
